# bots/tg_calculation_bot/tg_calculation_bot.service
[Service]
ExecStart=poetry run python /home/sergh/python-projects/bots/tg_calculation_bot/main.py

[Install]
WantedBy=multi-user.target